# debian/postinst
#! /bin/sh
# postinst script for smokeping
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the `postinst' is called with `abort-upgrade',
#     `abort-remove' or `abort-deconfigure'.

setup_permissions() {
    # Smokeping.cgi is run as user www-data, so we need it to be able
    # to write in /var/cache/smokeping/images and in /var/lib/smokeping/__cgi

    chown smokeping:smokeping /var/lib/smokeping
    chown www-data /var/cache/smokeping/images
    chown smokeping:www-data /var/lib/smokeping/__cgi
    chmod 2775 /var/lib/smokeping/__cgi
}

setup_apache() {
    # note that this is only done on initial installations
    if [ -d /etc/apache2/conf.d ] && \
       [ ! -e /etc/apache2/conf.d/smokeping ] && \
       which apache2ctl >/dev/null 2>&1
    then
        if apache2ctl configtest 2>/dev/null; then
            echo "enabling Apache configuration in /etc/apache2/conf.d/smokeping" 1>&2
            ln -s /etc/smokeping/apache2.conf /etc/apache2/conf.d/smokeping
            if apache2ctl configtest 2>/dev/null; then
                invoke-rc.d apache2 force-reload || true
            fi
        else
            echo "broken Apache configuration detected, not enabling /etc/smokeping/apache2.conf" 1>&2
        fi
    else
        echo "non-standard Apache configuration or missing apache2ctl, not enabling /etc/smokeping/apache2.conf" 1>&2
    fi
}


# OK, let's go
case "$1" in
    configure|reconfigure)

    getent passwd smokeping > /dev/null || \
      adduser --system --group --home /var/lib/smokeping --no-create-home --gecos "SmokePing daemon" smokeping
    if [ -z "$2" ]
    then
        setup_apache
    fi
    setup_permissions
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


